# devices/services/linearbot.service
[Unit]
Description=i2rt LinearRailVehicle controller (base + linear rail, 9 motors)
Documentation=https://github.com/i2rt-robotics/i2rt
# CAN comes up via devices/rules/flow_base.rules; wait for the interface before starting.
After=network.target sys-subsystem-net-devices-can0.device
Wants=sys-subsystem-net-devices-can0.device

[Service]
Type=simple
User=i2rt
WorkingDirectory=/srv/i2rt
# Wait (up to ~30s) for can0 to be present and up, in case udev brings it up after boot.
ExecStartPre=/bin/sh -c 'for i in $(seq 1 60); do ip link show can0 up >/dev/null 2>&1 && exit 0; sleep 0.5; done; echo "can0 not up after 30s" >&2; exit 1'
# On x86 (non-Raspberry Pi), append "--device <serial>" for the USB-GPIO converter; ignored on a Pi (native GPIO).
ExecStart=/srv/i2rt/.venv/bin/python /srv/i2rt/i2rt/flow_base/flow_base_controller.py --channel can0 --linear-rail

# The control loop requests SCHED_FIFO; grant the capability and RT priority so it succeeds.
AmbientCapabilities=CAP_SYS_NICE
LimitRTPRIO=99

Restart=on-failure
RestartSec=3
# Give the linear-rail homing/calibration time to finish on startup before health checks.
TimeoutStartSec=90

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
